# roles/famedly_sync/templates/famedly-sync.service.j2
[Unit]
Description=Sync users from LDAP to Zitadel

[Service]
Type=oneshot
ExecStart=/usr/bin/docker start -a famedly-sync
