# roles/famedly_sync/templates/famedly-sync.timer.j2
[Unit]
Description=Sync users from LDAP to Zitadel regularly

[Timer]
{% if nightly_sync %}
OnCalendar=*-*-* 11:45:00
RandomizeDelaySec=30m
{% else %}
OnBootSec=15min
OnUnitActiveSec=15min
{% endif %}

[Install]
WantedBy=timers.target
